# monitor/crontab
* * * * * /usr/local/bin/python /app/monitor.py
